# scripts/serverKtalyse.service
[Unit]
Description=Programe du serveur central de la société Ktalyse
After=network.target

[Service]
ExecStart=/usr/bin/node src/serverKtalyse.js
WorkingDirectory=/home/serveur_raspberry_rtls/src

Restart=always
User=pi

[Install]
WantedBy=multi-user.target
